# .config/systemd/user/hyprpaper.service
[Unit]
Description=Hyprpaper wallpaper daemon

[Service]
ExecStart=/usr/bin/hyprpaper
Restart=on-failure

[Install]
WantedBy=default.target
